# grocy_scanner/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:app]
command=python /app/app.py
directory=/app
